# neurocomputer/neurocomputer-watchdog.timer
[Unit]
Description=Run Neurocomputer Health Check Every 30 Seconds

[Timer]
OnBootSec=60
OnUnitActiveSec=30
AccuracySec=5

[Install]
WantedBy=timers.target
